# deploy/srtelnet.service
[Unit]
Description=second-reality-TELNET server (stream Second Reality over telnet)
Documentation=https://github.com/Paulie420/second-reality-TELNET
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=paulie420
Group=paulie420
WorkingDirectory=/home/paulie420/second-reality-TELNET
# Reference 30fps truecolor unit. For live management, prefer
# `tools/switch_fps.sh 20|30` which (re)writes this file and handles
# fps-specific flags automatically — see docs/deploy.md.
#
# The venv python has cap_net_bind_service set so it can bind low ports
# without needing sudo. Change --port to 23 if you want standard telnet.
ExecStart=/home/paulie420/second-reality-TELNET/.venv/bin/python -m srtelnet.server \
    --port 2323 \
    --frames /home/paulie420/second-reality-TELNET/frames-30fps \
    --fps 30 \
    --prewarm 80
Restart=on-failure
RestartSec=3
StandardOutput=journal
StandardError=journal

# Allow plenty of file descriptors for concurrent telnet connections.
# Each connection is ~1 fd; default soft limit (1024) would cap us.
LimitNOFILE=65536

# Minimal hardening — unprivileged LXCs can't set up mount/cgroup/kernel
# namespaces, so PrivateTmp/ProtectSystem/ProtectHome/Protect{Kernel,Cgroup}*
# all fail with status=226/NAMESPACE. NoNewPrivileges is a simple prctl flag
# with no namespace requirements, so we keep it.
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
